# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:synapse]
command=python /debug_synapse.py
directory=/data
autostart=true
autorestart=false
stdout_logfile=/var/log/supervisor/synapse.log
stderr_logfile=/var/log/supervisor/synapse-error.log

[program:localtunnel-sydent]
command=lt --port 8090 --subdomain chator-identity
autostart=true
autorestart=true
exitcodes=0
stdout_logfile=/var/log/localtunnel/sydent.log
stderr_logfile=/var/log/localtunnel/sydent-error.log

[program:livekit]
command=/usr/local/bin/livekit-server --port 7880 --keys "LKISupervisorSecret-Changer-in-production: chator_livekit_secret_change_me"
autostart=true
autorestart=true
stdout_logfile=/var/log/livekit/livekit.log
stderr_logfile=/var/log/livekit/livekit-error.log

[program:livekit-jwt]
command=/usr/local/bin/lk-jwt-service
directory=/data
autostart=true
autorestart=true
stdout_logfile=/var/log/livekit/jwt.log
stderr_logfile=/var/log/livekit/jwt-error.log
environment=LIVEKIT_URL="ws://localhost:7880",LIVEKIT_KEY="LKISupervisorSecret-Changer-in-production",LIVEKIT_SECRET="chator_livekit_secret_change_me",LIVEKIT_JWT_BIND=":8070"

# [program:coturn]
# NOTE: Coturn runs as a separate container (docker-compose.yml)
# It needs NET_ADMIN cap and runs as non-root — can't run inside main container
# Start with: docker compose up -d coturn

[program:sydent]
command=sydent
directory=/data
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/sydent.log
stderr_logfile=/var/log/supervisor/sydent-error.log

###############################################################################
# DNS-over-HTTPS proxy for Russia (bypasses DNS DPI blocking)
###############################################################################
[program:cloudflared-doh]
command=/usr/local/bin/cloudflared proxy-dns --port 5053 --address 127.0.0.1 --upstream https://1.1.1.1/dns-query --upstream https://dns.google/dns-query
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/cloudflared.log
stderr_logfile=/var/log/supervisor/cloudflared-error.log
# Respawn quickly if it crashes
exitcodes=0

###############################################################################
# Main nginx (port 8008: serves Element Web, Element Call, proxies Synapse)
###############################################################################
[program:element-web]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/element-web.log
stderr_logfile=/var/log/supervisor/element-web-error.log