# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:synapse]
command=python /debug_synapse.py
directory=/data
autostart=true
autorestart=false
stdout_logfile=/var/log/supervisor/synapse.log
stderr_logfile=/var/log/supervisor/synapse-error.log

[program:dex]
command=/usr/local/bin/dexidp serve /data/dex.yaml
directory=/data
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/dex.log
stderr_logfile=/var/log/supervisor/dex-error.log

[program:localtunnel-chator-server]
command=lt --port 8008 --subdomain chator-server
autostart=true
autorestart=true
exitcodes=0
stdout_logfile=/var/log/localtunnel/chator-server.log
stderr_logfile=/var/log/localtunnel/chator-server-error.log

[program:localtunnel-chator-dexter]
command=lt --port 5556 --subdomain chator-dexter
autostart=true
autorestart=true
exitcodes=0
stdout_logfile=/var/log/localtunnel/chator-dexter.log
stderr_logfile=/var/log/localtunnel/chator-dexter-error.log

[program:localtunnel-chator-cal]
command=lt --port 7880 --subdomain chator-kittimamumavabratamav
autostart=true
autorestart=true
exitcodes=0
stdout_logfile=/var/log/localtunnel/chator-cal.log
stderr_logfile=/var/log/localtunnel/chator-cal-error.log

[program:localtunnel-chator-cal-jwt]
command=lt --port 8070 --subdomain chator-iwi
autostart=true
autorestart=true
exitcodes=0
stdout_logfile=/var/log/localtunnel/chator-cal-jwt.log
stderr_logfile=/var/log/localtunnel/chator-cal-jwt-error.log

[program:localtunnel-sydent]
command=lt --port 8090 --subdomain chator-identity
autostart=true
autorestart=true
exitcodes=0
stdout_logfile=/var/log/localtunnel/sydent.log
stderr_logfile=/var/log/localtunnel/sydent-error.log

[program:livekit]
command=/usr/local/bin/livekit-server --port 7880 --keys "LKISupervisorSecret-Changer-in-production: chator_livekit_secret_change_me"
autostart=true
autorestart=true
stdout_logfile=/var/log/livekit/livekit.log
stderr_logfile=/var/log/livekit/livekit-error.log

[program:livekit-jwt]
command=/usr/local/bin/lk-jwt-service
directory=/data
autostart=true
autorestart=true
stdout_logfile=/var/log/livekit/jwt.log
stderr_logfile=/var/log/livekit/jwt-error.log
environment=LIVEKIT_URL="ws://localhost:7880",LIVEKIT_KEY="LKISupervisorSecret-Changer-in-production",LIVEKIT_SECRET="chator_livekit_secret_change_me",LIVEKIT_JWT_BIND=":8070"

# [program:coturn]
# NOTE: Coturn requires separate container or host setup due to user/group requirements
# Run separately: docker run -d --name coturn -p 3478:3478 coturn/coturn

[program:sydent]
command=sydent
directory=/data
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/sydent.log
stderr_logfile=/var/log/supervisor/sydent-error.log